# templates/meilisearch.service.j2
[Unit]
Description=Meilisearch
Documentation=https://www.meilisearch.com/docs
After=network.target network-online.target systemd-user-sessions.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{ meilisearch_home }}
EnvironmentFile=-{{ meilisearch_config_dir }}/meilisearch.env
ExecStart={{ meilisearch_binary_path }} --config-file-path {{ meilisearch_config_dir }}/meilisearch.toml
User={{ meilisearch_user }}
Group={{ meilisearch_group }}
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ meilisearch_home }} {{ meilisearch_config_dir }}

[Install]
WantedBy=multi-user.target
